# rootfs/etc/hotplug.d/iface/95-gl_ddns
#!/bin/sh

enable=`uci -q get gl_ddns.glddns.enabled`
if [ "$enable" = 1 -a "$INTERFACE" != "lan" ];then
	case "$ACTION" in
	ifdown |\
	ifup)
		/etc/init.d/gl_ddns restart
	;;
	esac
fi
