# rootfs/etc/hotplug.d/iface/99-gl-cloud
[ "$ACTION" = "ifup" -o "$ACTION" = "ifdown" ] || exit 0

ubus call gl-cloud notify-network-info &
