# rootfs/etc/hotplug.d/gl-tertf/00-gl_clients
ubus call gl-clients notify "{\"action\":\"$ACTION\", \"mac\":\"$MAC\"}"
