# rootfs/etc/hotplug.d/dhcp/00-gl_dhcp_ping
#!/bin/sh

[ -n "$IPADDR" ] || exit 0

if [ "$ACTION" = "update" -o "$ACTION" = "add" ]; then
    ping $IPADDR -q -w 5 > /dev/null 2>&1 &
fi
