# rootfs/usr/lib/opkg/info/alist.list
/etc/config/alist
/etc/alist/data.db
/etc/init.d/alist
/lib/upgrade/keep.d/alist
/usr/bin/alist
