# rootfs/usr/lib/opkg/info/gl-oui-rpc.list
/usr/lib/lua/oui/rpc.lua
/usr/lib/lua/oui/db.lua
/usr/lib/oui-httpd/rpc/iwinfo
/usr/sbin/gl-ngx-session
/etc/init.d/gl-ngx-session
/lib/upgrade/keep.d/gl-oui-rpc
/usr/share/gl-ngx/oui-access.lua
/usr/share/gl-ngx/oui-rpc.lua
/usr/share/gl-ngx/oui-init.lua
/usr/lib/oui-httpd/rpc/ubus
/usr/share/gl-ngx/oui-upload.lua
/etc/nginx/gl-conf.d/access.conf
/www/cgi-bin/glc
/usr/lib/lua/resty/http.lua
/usr/lib/lua/oui/iwinfo.lua
/etc/config/oui-httpd
/usr/lib/lua/resty/http_headers.lua
/usr/lib/lua/resty/upload.lua
/usr/lib/lua/oui/system.so
/etc/nginx/oui_nginx.conf
/etc/nginx/conf.d/gl.conf
/usr/lib/oui-httpd/rpc/uci
/usr/lib/lua/oui/utils.lua
/usr/lib/oui-httpd/rpc/ui
/etc/uci-defaults/80_nginx-oui
/etc/init.d/fcgiwrap
/usr/bin/fcgiwrap
/usr/lib/lua/oui/utils/utils.so
/usr/share/gl-ngx/oui-download.lua
/usr/lib/lua/resty/http_connect.lua
/usr/lib/lua/oui/fs.so
/usr/lib/lua/oui/network.so
/usr/lib/lua/resty/resolver.lua
/usr/lib/lua/oui/ubus.lua
