# rootfs/usr/lib/opkg/info/gl-sdk4-tethering.list
/etc/init.d/gl_tethering
/etc/hotplug.d/usb/31-tethering
/usr/lib/oui-httpd/rpc/tethering
