# rootfs/usr/lib/opkg/info/gl-sdk4-2b.list
/etc/init.d/gl-update-logo
/usr/bin/gl-update-logo
/lib/preinit/82_indicate_2b
/etc/uci-defaults/z-2b-cfg
/lib/preinit/81_merge_2b_cfg
