# rootfs/usr/lib/opkg/info/gl-sdk4-modem.list
/rootfs-overlay/usr/bin/modem_upgrade_check
/etc/chatscripts/3gipv6.chat
/usr/bin/qcm
/usr/bin/sms_manager
/rootfs-overlay/usr/bin/verizon_apn.sh
/etc/config/custom_apn
/rootfs-overlay/lib/netifd/proto/3g.sh
/usr/bin/modem_AT
/usr/bin/gl_modem
/lib/upgrade/keep.d/gl-sdk4-modem
/rootfs-overlay/usr/bin/switch_sim_slot
/etc/config/apnprofile
/usr/bin/quec_upload
/etc/uci-defaults/99-set-apnprofile
/etc/init.d/modem
/usr/bin/QFirehose
/etc/init.d/modem-init
/etc/hotplug.d/iface/17-dual-sim
/rootfs-overlay/usr/bin/get_modem_bus.lua
/etc/init.d/sms_manager
/rootfs-overlay/lib/netifd/proto/callip.sh
/rootfs-overlay/lib/netifd/proto/qmi.sh
/etc/uci-defaults/99-set-monitor
/etc/hotplug.d/usb/30_modem.sh
/www/js/apn-config.json.gz
/etc/uci-defaults/99-verizon
/rootfs-overlay/lib/netifd/proto/qcm.sh
/etc/gl_crontabs/crontabs.d/modem_status_monitor
/etc/chatscripts/evdoipv6.chat
/lib/functions/modem.sh
/etc/config/glmodem
/etc/verizon.type
/usr/lib/oui-httpd/rpc/modem.so
/rootfs-overlay/usr/bin/set_modem_ttl
/rootfs-overlay/usr/bin/modem_manager.lua
/etc/init.d/modem_auto
/rootfs-overlay/lib/netifd/proto/directip.sh
/etc/hotplug.d/usb/40_recover_modem.sh
/rootfs-overlay/usr/bin/sendsms
