# rootfs/usr/lib/opkg/info/gl-sdk4-led.list
/usr/lib/oui-httpd/rpc/led
/etc/init.d/gl_led
/usr/bin/gl_i2c_led
/usr/bin/gl_led
/etc/config/gl_led
