# rootfs/usr/lib/opkg/info/gl-sdk4-repeater.list
/etc/gl-reset-network.d/reset-repeater
/etc/init.d/repeater
/usr/lib/oui-httpd/rpc/repeater
/usr/lib/lua/wifi-event.so
/etc/config/repeater
/usr/sbin/repeater
