# rootfs/usr/lib/opkg/info/gl-sdk4-s2s.list
/etc/uci-defaults/gl-s2s.defaults
/etc/config/gl_s2s
/usr/lib/oui-httpd/rpc/s2s.so
/usr/bin/echo_server
/etc/init.d/gl_s2s
