# rootfs/usr/lib/opkg/info/gl-sdk4-edgerouter.list
/etc/edgerouter.fw
/etc/uci-defaults/14-fix-edgerouter-keepsetting
/lib/functions/drop-in-function.sh
/etc/init.d/edgerouter
/usr/lib/oui-httpd/rpc/edgerouter
/etc/config/edgerouter
