# rootfs/usr/lib/opkg/info/gl-sdk4-ipv6.list
/etc/config/glipv6
/etc/init.d/gl_ipv6
/etc/uci-defaults/13_add_nat6_to_firewall
/etc/firewall.nat6
/etc/hotplug.d/dhcpv6/01-gl_relay_add_neigh_proxy
/etc/hotplug.d/dhcpv6/00-gl_relay_add_br-lan_route
/usr/bin/get_wan_interface
/usr/bin/sipcalc
/etc/hotplug.d/iface/99-guest-laddr
/usr/lib/oui-httpd/rpc/ipv6
/etc/uci-defaults/13_fix_ipv6_config
