# rootfs/usr/lib/opkg/info/gl-sdk4-timer.list
/usr/bin/gl_timer_control_wifi
/etc/init.d/gl_timer
/etc/config/gl_timer
/usr/bin/gl_timer_control_disk
/usr/lib/oui-httpd/rpc/timer
