# rootfs/usr/lib/opkg/info/gl-sdk4-base-files.list
/rootfs-overlay/etc/uci-defaults/zzzz_gl_init_done
/rootfs-overlay/etc/uci-defaults/99-firewall
/rootfs-overlay/etc/uci-defaults/network_gl_01_disable_classlessroute
/rootfs-overlay/etc/hotplug.d/iface/20-firewall
/rootfs-overlay/sbin/sysupgrade
/rootfs-overlay/lib/preinit/82_fix_i2c_gpio_custom_preinit
/rootfs-overlay/etc/rc.button/switch
/rootfs-overlay/etc/uci-defaults/network_gl
/rootfs-overlay/usr/bin/modem_signal
/rootfs-overlay/etc/uci-defaults/99-set_modem.sh
/rootfs-overlay/lib/netifd/ppp6-up
/rootfs-overlay/etc/uci-defaults/15_fix_glconfig_track_ip
/rootfs-overlay/etc/uci-defaults/19-set-hardware-params
/rootfs-overlay/etc/hotplug.d/block/10-mount
/rootfs-overlay/lib/functions/gl_util.sh
/rootfs-overlay/etc/init.d/mtk-hwnat-post
/rootfs-overlay/etc/uci-defaults/02_gl_fan
/rootfs-overlay/etc/board.d/02_gl_system
/rootfs-overlay/usr/bin/clean_client_conntrack
/rootfs-overlay/etc/uci-defaults/dhcp_server
/rootfs-overlay/etc/init.d/modem_signal
/rootfs-overlay/etc/config/board_special
/rootfs-overlay/usr/bin/get_interface_status.lua
/rootfs-overlay/etc/uci-defaults/01_gl_system
/rootfs-overlay/etc/uci-defaults/17_fix_glconfig_wan2lan_init
/rootfs-overlay/etc/gl-switch.btn
/rootfs-overlay/etc/uci-defaults/16_fix_ovpn_dns_leak_keep
/rootfs-overlay/etc/config/glconfig
/rootfs-overlay/etc/rc.local
/rootfs-overlay/etc/uci-defaults/17_fix_ax1800_upgrade_url
/rootfs-overlay/lib/preinit/83_mini_free_kbytes
/rootfs-overlay/etc/rc.button/reset
/rootfs-overlay/etc/config/switch-button
/rootfs-overlay/etc/board.d/03_gl_network
/rootfs-overlay/etc/init.d/gl_init
/rootfs-overlay/lib/preinit/82_gl_fw_type
/rootfs-overlay/etc/uci-defaults/03_gl_switch
