# rootfs/usr/lib/opkg/info/gl-sdk4-plugins.list
/usr/bin/update_plugins_db
/etc/config/plugins
/usr/lib/oui-httpd/rpc/plugins.so
/etc/init.d/plugins
